# scripts/clash.cron
#!/system/bin/sh 
scripts=`realpath $0`
scripts_dir=`dirname ${scripts}`

step=5 # interval in seconds, cannot be greater than 60  
i=0
while ((i < 60)) ; do  
	i=$i+step
    ${scripts_dir}/clash.tool -m # > /data/adb/clash/run/cron.log
    sleep $step
done  
exit 0